# supervisor.conf
[supervisord]
nodaemon=true

[program:minio]
command=/usr/local/bin/minio server /local_storage --address 0.0.0.0:9000 --console-address 0.0.0.0:9090
autostart=true
autorestart=true
stderr_logfile=/var/log/minio.err.log
stdout_logfile=/var/log/minio.out.log
user=minio-user

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
user=root